# essentials/.zsh/profile.d/01-basics
#################
# The very basics
#################

export EDITOR=vi
export VISUAL=vi
export PAGER=less

#############
# Web browser
#############
if [ -n "$DISPLAY" ]; then
    BROWSER=chromium
else
    BROWSER=links
fi
export BROWSER

###################
# Terminal settings
###################
if [[ -z "$TMUX" ]] ;then
    export TERM=xterm-256color
fi
export TERMINFO=${DOT_PATH}/share/terminfo
