# i3/.zsh/profile.d/13-i3
# vim: ft=zsh
########################################
# Settings for i3wm
########################################

# i3 executable location
I3WM_BIN=/usr/bin/i3

# Log file location
I3WM_LOG_FILE=~/.i3/i3log-$(date +'%F-%k-%M-%S')
