# zsh/.zsh/profile.d/10-dot
# Path to dotfiles folder
export DOT_HOME=$HOME/.dotfiles

# Custom installs PATH
export DOT_PATH=$HOME/.local/dot

# LD_LIBRARY_PATH is used in preference to any run time or default system linker path.
export LD_LIBRARY_PATH=$DOT_PATH/lib

# Standard PATH at the user level
export PATH=$DOT_PATH/bin:$PATH
