# zsh/.zsh/profile.d/10-dot
#################################
# Standard PATH at the user level
#################################
PATH=$DOT_PATH/bin:$PATH
