# zsh/.zsh/profile.d/90-dotenv
# vim: ft=zsh
############################################
# Fill the environment with
# custom variables via dot-env
############################################
dot-env > $DOT_ENV_OUT 2> /dev/null
[ -f $DOT_ENV_OUT ] && . $DOT_ENV_OUT
