# zsh/.zsh/profile.d/91-custom-path
# vim: ft=zsh
################################
# Settings for custom PATH
################################
if [[ ! -z "$DOT_CUSTOM_PATH" ]]; then
    PATH=$DOT_CUSTOM_PATH:$PATH
fi
