# zsh/.zsh/profile.d/99-dotenv
############################################
# Fill the environment with
# custom variables via dot-env
############################################
eval $(dot-env 2> /dev/null)
