# X11/.zsh/profile.d/12-x11
# vim: ft=zsh
########################################
# Settings for X11
########################################
export DISPLAY_DPI=96 # Universal value for setting up dpi
