# sync/.zsh/profile.d/15-sync
#############################################
# file sync root directory for automatic stowing
#############################################
export SYNC_HOME=${HOME}/.dotsync
