# sync/.zsh/profile.d/15-sync
#############################################
# btsync root directory for automatic stowing
#############################################
export BTSYNC_HOME=${HOME}/.dotsync
