# sync/.zsh/profile.d/15-sync
#############################################
# List of directories for automatic stowing
# via fsync
#############################################
export SYNC_DIRS=${HOME}/.dotsync
