# sync/.zsh/profile.d/15-stow-dirs
# vim: ft=zsh
#############################################
# Stow directories are directories to be
# treated as symlink farms managed by stow
# just as DOT_HOME with the only difference
# that these are watched and automatically
# managed by their own parent swatch process
# (see swatch for more details)
#############################################
export STOW_DIRS="${HOME}/.dotsync"
