# android/.zsh/profile.d/93-android-sdk
# vim: ft=zsh
################################
# Settings for Android SDK:
################################

[[ ! -z "$ANDROID_HOME" ]] && PATH=$ANDROID_HOME/tools:$ANDROID_HOME/platform-tools:$PATH
