# base16/.zsh/profile.d/11-base16
# Set default scheme to be used by base16
export BASE16_SCHEME=behelit.yml
