# golang/.zsh/profile.d/21-golang
# Settings for golang
export GOROOT=~/.local/dot/golang
export GOPATH=~/Projects/golang
export PATH=$GOROOT/bin:$GOPATH/bin:$PATH
