# golang/.zsh/profile.d/21-golang
# Settings for golang
export GOPATH=~/Projects/golang
export PATH=$GOPATH/bin:$PATH
