# golang/.zsh/profile.d/21-golang
#####################
# Settings for golang
#####################
export GOROOT=~/.local/dot/golang
export GOPATH=~/Projects/golang

# so golang-made programs can be seen
export PATH=$GOROOT/bin:$GOPATH/bin:$PATH
